# systemd/crimson-osd@.service.in
[Unit]
Description=Crimson object storage daemon osd.%i
After=network-online.target local-fs.target time-sync.target
Wants=network-online.target local-fs.target time-sync.target

[Service]
LimitNOFILE=infinity
ExecStart=/var/lib/ceph/bin/crimson-osd --smp 1 -i %i -c /var/lib/ceph/ceph.conf

Restart=on-failure
StartLimitBurst=3
RestartSec=1

[Install]
WantedBy=multi-user.target
